# etc/systemd/system/zammad_spam_rspamd.timer
[Unit]
Description=Periodically send Tickets tagged as spam to rspamd

[Timer]
OnCalendar=*:0/30
Persistent=true

[Install]
WantedBy=timers.target
